# asneg-db-server/OpcUaDB/src/OpcUaDB/linux/usr/lib/systemd/system/OpcUaDB.service.in
[Unit]
Description=OPC UA DB server

[Service]
Environment="LD_LIBRARY_PATH=/usr/lib/OpcUaStack:${LD_LIBRARY_PATH}"
ExecStart=/usr/bin/OpcUaServer@OPCUASTACK_VERSION_MAJOR@ /etc/OpcUaStack/OpcUaDB/OpcUaServer.xml
Restart=on-failure

[Install]
WantedBy=multi-user.target
